# app/deploy/profit-admin-api.service
[Unit]
Description=Outscore Profit Admin API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/agents/outscore_profit
Environment="PATH=/opt/agents/outscore_profit/venv/bin"
EnvironmentFile=/opt/agents/outscore_profit/.env
ExecStart=/opt/agents/outscore_profit/venv/bin/uvicorn profit_api.app:app --host 127.0.0.1 --port 8010
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

